# dude-agent.service
[Unit]
Description=Dude Self-Improving AI Agent
After=network.target

[Service]
Type=simple
User=john
WorkingDirectory=/home/john/dude
ExecStartPre=/usr/bin/bash -c "[ -z \"$GEMINI_JSON_TOKEN\" ] || echo \"$GEMINI_JSON_TOKEN\" > /home/john/.pi/agent/auth.json"
ExecStart=/usr/bin/nix develop --command npm start
Restart=always
EnvironmentFile=/home/john/dude/.env

[Install]
WantedBy=multi-user.target
